# packaging/systemd/hambridge.service
[Unit]
Description=HaMBridge: Hardware-MQTT Bridge (evdev to MQTT)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hambridge
Group=hambridge
# Config path — adjust if your distro uses a different layout.
ExecStart=/usr/bin/hambridge --config /etc/hambridge/config/hambridge.yaml
Restart=on-failure
RestartSec=5

# Baseline hardening (compatible with reading /etc and publishing MQTT)
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
